# config/autostart/Conky.desktop
[Desktop Entry]
Type=Application
Name=Conky
Exec=sh -c "sleep 10 && conky"
StartupNotify=false
Terminal=false
Hidden=false
